# deploy/video_processor.service
[Unit]
Description=Audio Processing Service
After=network-online.target
StartLimitIntervalSec=0

[Service]
Group=www-data
Type=simple
Restart=always
RestartSec=1
LimitRTPRIO=99
LimitRTTIME=infinity
WorkingDirectory=/var/lib/chemistry-dashboard/
ExecStart= /var/lib/chemistry-dashboard/video_processing/venv/bin/python3 video_processing/server.py

[Install]
WantedBy=multi-user.target
